# conf/xcalar.slice
[Unit]
Description=Xcalar Slice
Documentation=https://xcalar.com/documentation/Content/DocPortalHome.htm
Before=slices.target
