# pkg/docker/supervisor.conf
[supervisord]
nodaemon=true
minfds=14096
minprocs=1800
loglevel=info
pidfile=%(ENV_XCE_WORKDIR)s/supervisord.pid
directory=%(ENV_XCE_WORKDIR)s
nocleanup=false
childlogdir=%(ENV_XCE_LOGDIR)s
logfile=%(ENV_XCE_LOGDIR)s/supervisor.log
logfile_maxbytes=100MB
#user=xcalar
#group=xcalar
#environment=USER="xcalar",HOME="/home/xcalar",LOGNAME="xcalar"


[supervisorctl]
serverurl=unix://%(ENV_XCE_WORKDIR)s/supervisor.sock

[unix_http_server]
file=%(ENV_XCE_WORKDIR)s/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[group:xcalar]
programs=expserver,xcmgmtd,caddy,sqldf,xcmonitor

[program:expserver]
environment=NODE_PATH="/usr/lib/node_modules"
command=node %(ENV_XLRGUIDIR)s/services/expServer/expServer.js
directory=%(ENV_XLRGUIDIR)s/services/expServer
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=%(ENV_XCE_LOGDIR)s/%(program_name)s.out
stdout_logfile_maxbytes=100MB

[program:xcmgmtd]
command=%(ENV_XLRDIR)s/bin/xcmgmtd "%(ENV_XCE_CONFIG)s"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=%(ENV_XCE_LOGDIR)s/%(program_name)s.out
stdout_logfile_maxbytes=100MB

[program:caddy]
startsecs=10
startretries=5
command=%(ENV_XLRDIR)s/bin/caddy -quiet -conf "%(ENV_XCE_CADDYFILE)s"
autostart=true
autorestart=unexpected
redirect_stderr=true
stdout_logfile=%(ENV_XCE_LOGDIR)s/%(program_name)s.out
stdout_logfile_maxbytes=100MB

[program:jupyter]
environment=HOME=/home/xcalar,PYTHONUSERBASE="%(ENV_HOME)s/.local",HOME=%(ENV_XCE_HOME)s
command=%(ENV_XLRDIR)s/bin/python3.6 %(ENV_XLRDIR)s/bin/jupyter-notebook
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=%(ENV_XCE_LOGDIR)s/%(program_name)s.out
stdout_logfile_maxbytes=100MB

[program:sqldf]
command=java -jar %(ENV_XLRDIR)s/lib/xcalar-sqldf.jar -j
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=%(ENV_XCE_LOGDIR)s/%(program_name)s.out
stdout_logfile_maxbytes=100MB

[program:xcmonitor]
command=%(ENV_XLRDIR)s/bin/xcmonitor -n %(ENV_XCE_NODEID)s -m %(ENV_XCE_NUMNODES)s -c %(ENV_XCE_CONFIG)s
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=%(ENV_XCE_LOGDIR)s/%(program_name)s.out
stdout_logfile_maxbytes=100MB
